# development/build/supervisor/test/run.conf
[program:phototrend]
directory=/data/app/phototrend/
command=/root/anaconda3/envs/phototrend/bin/python app.py --xformers --setup_mode --port 7777
autostart=true
numprocs=1
startsecs=5
autorestart=true
startretries=10
user=root
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=20
process_name=%(process_num)s
stdout_logfile=/data/log/phototrend/%(group_name)s.log
stderr_logfile=/data/log/phototrend/error_%(group_name)s.log